# data/99-wpc-racecam-leds.rules
ACTION=="add", KERNEL=="LED1", SUBSYSTEM=="leds", RUN+="/bin/chown -R root:racecam /sys%p", RUN+="/bin/chmod -R g+rw /sys%p"
